# service/mysql.pull.service
[Unit]
Description=mysql-pull

[Service]
ExecStart=/usr/bin/docker pull mysql
